# tools/com.chiero.syncnote.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  chiero.jp の「記録」を note に合わせて日次で更新する。

  設置（chiero.jp を新サイトへ切り替えた後に行う。それまでは意味がない）:
    cp tools/com.chiero.syncnote.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.chiero.syncnote.plist
  確認:
    launchctl list | grep chiero
    tail ~/Desktop/Claude/chiero_site/tools/sync_note.log
  外す:
    launchctl unload ~/Library/LaunchAgents/com.chiero.syncnote.plist

  時刻は 07:10。既存の窓の検知網（07:00 collect）とずらしてある。
  Mac が寝ていた場合は起きた後に一度だけ走る（RunAtLoad + StartCalendarInterval）。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chiero.syncnote</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>$HOME/Desktop/Claude/chiero_site/tools/sync_note_daily.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>7</integer>
    <key>Minute</key><integer>10</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/chiero_syncnote.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/chiero_syncnote.err</string>
</dict>
</plist>
